# root/etc/s6-overlay/s6-rc.d/audiobookshelf/run
#!/usr/bin/with-contenv /bin/bash
# shellcheck shell=bash

echo "Starting audiobookshelf..."

# Set environment variables (LinuxServer.io baseimage already runs as user 'abc')
PORT=${PORT:-80}
HOST=${HOST:-0.0.0.0}
CONFIG_PATH=${CONFIG_PATH:-/config}
METADATA_PATH=${METADATA_PATH:-/config/metadata}
FFMPEG_PATH=${FFMPEG_PATH:-/usr/bin/ffmpeg}
FFPROBE_PATH=${FFPROBE_PATH:-/usr/bin/ffprobe}

echo "  Port: ${PORT}"
echo "  Host: ${HOST}"
echo "  Config: ${CONFIG_PATH}"
echo "  Metadata: ${METADATA_PATH}"

# Change to application directory and start
# Note: No s6-setuidgid needed - LinuxServer.io baseimage already runs as user 'abc'
cd /app || exit 1

exec node index.js