# root/etc/s6-overlay/s6-rc.d/init-audiobookshelf-config/up
/etc/s6-overlay/s6-rc.d/init-audiobookshelf-config/run